# deployment/systemd/micro-ros-agent.service
[Unit]
Description=Micro-ROS Agent for PX4 Communication
After=network.target

[Service]
Type=simple
User=khadas
Group=khadas
WorkingDirectory=%h/humiro_fire_suppression
Environment="HOME=%h"
Environment="ROS_DOMAIN_ID=0"
Environment="ROS_NAMESPACE=drone1"
Environment="RMW_IMPLEMENTATION=rmw_fastrtps_cpp"
ExecStart=%h/humiro_fire_suppression/scripts/runtime/start_micro_ros_agent_wrapper.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
